# systemd/solarmon.service
[Unit]
Description=Huawei SUN2000 monitor
After=network-online.target syslog.target

[Service]
Type=simple
WorkingDirectory=/home/pi/solarmon
ExecStart=/usr/bin/python3 /home/pi/solarmon/main.py
StandardOutput=syslog
StandardError=syslog
Restart=on-failure

[Install]
WantedBy=multi-user.target
